# deploy_confs/supervisor/gbm_client.conf
[program:gbm_client]
command = /snap/bin/dotnet /webapps/gbmseg_dotnet/deploy/gbm_client/GBMWeb.Client.dll ; Command to start app
user = gbmseg                                                          ; User to run as
stdout_logfile = /webapps/gbmseg_dotnet/logs/gbm_client_supervisor.log   ; Where to write log messages
redirect_stderr = true                                                     ; Save stderr in the same log
;environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                           ; Set UTF-8 as default encoding
